# hermes/linux/odroid/systemd/hermes-dashboard-watchdog.timer
[Unit]
Description=Run HERMES dashboard health watchdog every 10 seconds

[Timer]
OnBootSec=20
OnUnitActiveSec=10
AccuracySec=1
Unit=hermes-dashboard-watchdog.service

[Install]
WantedBy=timers.target
