# hermes/linux/odroid/systemd/hermes-retention.service
[Unit]
Description=HERMES SQLite time-series retention and downsampling (oneshot)
After=network.target

[Service]
Type=oneshot
User=odroid
Environment=HOME=/home/odroid
Environment=HERMES_DB=/home/odroid/hermes-data/db/hermes.sqlite3
ExecStart=python3 /home/odroid/hermes-src/hermes/linux/odroid/scripts/retention.py
TimeoutStartSec=0
StandardOutput=journal
StandardError=journal
