# hermes/linux/odroid/systemd/hermes-oled-context.timer
[Unit]
Description=Push OLED context about every 30s

[Timer]
OnBootSec=20
OnUnitActiveSec=30
AccuracySec=5
Unit=hermes-oled-context.service

[Install]
WantedBy=timers.target
