# hermes/linux/odroid/systemd/hermes-oled-context.service
[Unit]
Description=Push OLED context to nRF (oneshot, non-overlapping)
After=network.target

[Service]
Type=oneshot
User=odroid
Environment=HOME=/home/odroid
WorkingDirectory=/home/odroid/hermes-src/hermes
# Single flight: skip this tick if the previous run is still holding the lock.
ExecStart=/usr/bin/flock -n /tmp/hermes-oled-context.lock /home/odroid/hermes-src/hermes/tools/push_oled_context.sh
TimeoutStartSec=90
StandardOutput=journal
StandardError=journal
