# hermes/linux/odroid/systemd/hermes-watchdog.timer
[Unit]
Description=Run HERMES stale stream watchdog every 10 seconds

[Timer]
OnBootSec=15
OnUnitActiveSec=10
AccuracySec=1
Unit=hermes-watchdog.service

[Install]
WantedBy=timers.target
