# hermes/linux/odroid/systemd/hermes-retention.timer
[Unit]
Description=Run HERMES DB retention daily at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=hermes-retention.service

[Install]
WantedBy=timers.target
