# hermes/linux/odroid/systemd/hermes-oled-context.timer
[Unit]
Description=Push OLED context every 10s

[Timer]
OnBootSec=10
OnUnitActiveSec=10
AccuracySec=1
Unit=hermes-oled-context.service

[Install]
WantedBy=timers.target
